# .obsidian_sync.service
[Unit]
Description=Backup obsidian notes

[Service]
Type=simple
ExecStart=/usr/local/bin/obsidian_backup.sh
